# fino-cdm-softnode/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh_server]
command=service ssh start
redirect_stderr=false

[program:cloudera_agent_server]
command=service cloudera-scm-agent start
redirect_stderr=false