# fino-cdm-repo/supervisord.conf
[supervisord]
nodaemon=true

[program:apache_server]
command=/etc/init.d/apache2 restart
redirect_stderr=true
stdout_logfile=/var/log/apache2/apache2.log